# conf/upstart/xvfb.conf
# Xvfb
#
# Runs Xvfb so that unoconv can work

description     "Xvfb"

start on (filesystem and started dbus)
stop on (stopping dbus)

respawn

script
    Xvfb :1
end script
