# conf/dyndnschecker/debian/nublic-dyndnschecker.cron.d
*/5 * * * * root /usr/bin/dyndnschecker
